# backend/proxy-pool/supervisord.conf
; supervisord -- lives inside the litetube-3proxy container.
;
; 3proxy is run as the 'litetube' user, with its cfg at /var/lib/3proxy/3proxy.cfg.
; The FastAPI app SIGUPs the 3proxy parent on every user_add / user_remove
; (via the cfg rendered by proxy_3proxy.render_cfg()) so we do not need to
; stop+start the program on user changes -- existing connections keep flowing.

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
minfds=1024
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:3proxy]
command=/usr/bin/3proxy /var/lib/3proxy/3proxy.cfg
directory=/var/lib/3proxy
autostart=true
autorestart=true
startretries=10
startsecs=2
priority=10
user=litetube
stdout_logfile=/var/log/3proxy/supervisor-stdout.log
stderr_logfile=/var/log/3proxy/supervisor-stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
environment=PRIORITY="10"

[program:healthcheck-curl]
; Lightweight watchdog: tries to TCP-connect to 11001/11002 every 60s.
; If both probes fail, the line in stderr_logfile triggers fail2ban or
; off-host Prometheus blackbox-exporter. Useful as a coarse "is the
; process dead?" alarm independent of the FastAPI health-checker.
command=/bin/sh -c 'while true; do curl --connect-timeout 5 --max-time 8 telnet://127.0.0.1:11001 >/dev/null 2>&1; curl --connect-timeout 5 --max-time 8 telnet://127.0.0.1:11002 >/dev/null 2>&1; sleep 60; done'
autostart=true
autorestart=false
priority=20
user=litetube
stdout_logfile=/var/log/3proxy/healthcheck-curl.log
stderr_logfile=/var/log/3proxy/healthcheck-curl.log
